# vault.d/etc/systemd/system/vault-init.timer
[Unit]
Description=Ensure Hashicorp Vault is initialised

[Timer]
OnBootSec=3600
OnUnitActiveSec=3600

[Install]
WantedBy=timers.target
