# matchbox.d/etc/systemd/system/matchbox-init.timer
[Unit]
Description=Ensure CoreOS matchbox init

[Timer]
OnBootSec=86400
OnUnitActiveSec=86400

[Install]
WantedBy=timers.target
