# consul.d/etc/systemd/system/consul-init.timer
[Unit]
Description=Ensure Hashicorp Consul init

[Timer]
OnBootSec=86400
OnUnitActiveSec=86400

[Install]
WantedBy=timers.target
